# templates/systemd/slurmd.service.erb
[Unit]
Description=Slurm node daemon
After=network.target
ConditionPathExists=<%= scope.lookupvar('slurm::slurm_conf_path') %>

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/slurmd
ExecStart=/usr/sbin/slurmd $SLURMD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=<%= scope.lookupvar('slurm::pid_dir') %>/slurmd.pid
KillMode=process
LimitNOFILE=51200
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target
