# templates/systemd/slurmdbd.service.erb
[Unit]
Description=Slurm DBD accounting daemon
After=network.target
ConditionPathExists=<%= scope.lookupvar('slurm::slurmdbd_conf_path') %>

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/slurmdbd
ExecStart=/usr/sbin/slurmdbd $SLURMDBD_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=<%= scope.lookupvar('slurm::pid_dir') %>/slurmdbd.pid

[Install]
WantedBy=multi-user.target
